# dot_config/systemd/user/ssh-agent.service
[Unit]
Description=SSH key agent

[Service]
Type=simple
ExecStart=/home/daniel/.local/bin/my-ssh-agent

[Install]
WantedBy=default.target
